# scann3r-supervisor.conf
[program:scann3r]
command=node src/index.js
directory=/home/pi/scann3r/
numprocs=1
user=pi
stdout_logfile=/var/log/supervisor/scann3r.log
autorestart=true


